# internetRadio-1.3.0-wifi-manager/services/radiomonitor.service
[Unit]
Description=Internet Radio System Monitor
After=internetradio.service NetworkManager.service
Wants=internetradio.service NetworkManager.service

[Service]
Type=simple
User=radio
Group=radio
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/radio/.Xauthority
WorkingDirectory=/home/radio/internetRadio
ExecStart=/usr/bin/lxterminal --geometry=100x30 --title="Radio Monitor" -e "/usr/bin/python3 -c 'from src.utils.system_monitor import SystemMonitor; SystemMonitor().run()'"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 